# bin/cron/clearAllSessions.cron
52 22 * 3,6,9,12 0 mysql -u full_publish full_production --password < /home/full/web/production/bin/sql/clearSessions.sql >> /home/full/web/logs/clearSession_cron.txt; cd /home/full/web/production'; ./bin/shell/clearcache2.sh;

52 22 * 3,6,9,12 0 mysql -u full_publish full_ezdev --password < /home/full/web/dev/public_html/bin/sql/clearSessions.sql >> /home/full/web/logs/clearSession_cron.txt; cd /home/full/web/dev/public_html/; ./bin/shell/clearcache2.sh;
