# bin/cron/export_Froogle.cron
52 1,2,3,4,5,6,7,8,9,23 * * * curl http://dev.example.com/trade/export >> /home/full/web/logs/froogle_cron.txt;
